# roles/wireguard_shaper/templates/wireguard-shaper.service.j2
[Unit]
Description=WireGuard Shaper

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wireguard-shaper.sh